# etc/cron.d/orchestrator
# etc/cron.d/simo-orchestrator
# To be deployed at /etc/cron.d/simo-orchestrator
NIX_BIN=/usr/local/simox/result/bin

# Run your custom 5-minute jobs safely isolated
# */5 * * * * root run-parts /etc/simo-cron.5min # not apt
*/5 * * * * root /bin/sh -c 'export PATH="$NIX_BIN:$PATH"; exec bash -c "run-parts /etc/simo-cron.5min"'

# Run your custom hourly jobs safely isolated
0 * * * * root /bin/sh -c 'export PATH="$NIX_BIN:$PATH"; exec bash -c "run-parts /etc/simo-cron.hourly"'

# Run your custom monthly jobs safely isolated
0 0 1 * * root /bin/sh -c 'export PATH="$NIX_BIN:$PATH"; exec bash -c "run-parts /etc/simo-cron.monthly"'
